# zcfan-resume.service
[Unit]
Description=Reload zcfan after waking up from suspend
After=sysinit.target
After=suspend.target
After=suspend-then-hibernate.target
After=hybrid-sleep.target
After=hibernate.target

[Service]
Type=oneshot
ExecStart=/usr/bin/pkill -x -USR2 zcfan

[Install]
WantedBy=sleep.target
